# raspiclient/rpiclient.service
[Unit]
Description = Raspberry Pi Client service
Requires = influxdb.service openhab2.service
After = influxdb.service openhab2.service

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/usr/bin/java -jar /home/pi/raspiclient-0.0.1-jar-with-dependencies.jar

[Install]
WantedBy=multi-user.target